# ssot/agent_runbooks/assets/acer_watchdog.timer
[Unit]
Description=Doraemon Acer watchdog timer

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
AccuracySec=15s
Unit=doraemon-acer-watchdog.service

[Install]
WantedBy=timers.target

